# scripts/com.qc-dashboard.daily-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qc-dashboard.daily-push</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/laitianyou/WorkBuddy/20260408163242/qc-dashboard/scripts/tao_daily_push_launchd.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/laitianyou/WorkBuddy/20260408163242/qc-dashboard</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>ToleratedDelay</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/laitianyou/WorkBuddy/20260408163242/qc-dashboard/data/logs/launchd-daily-push.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/laitianyou/WorkBuddy/20260408163242/qc-dashboard/data/logs/launchd-daily-push.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/laitianyou/.npm-global/bin:/Library/Frameworks/Python.framework/Versions/3.14/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/Users/laitianyou/.local/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/laitianyou/WorkBuddy/20260408163242/qc-dashboard/src</string>
        <key>HOME</key>
        <string>/Users/laitianyou</string>
        <key>WECOM_WEBHOOK_KEY</key>
        <string>c8bed28c-bd5c-499a-b460-0e62a7376d4e</string>
    </dict>

    <key>UserName</key>
    <string>laitianyou</string>
</dict>
</plist>
